# emperor.conf
# Emperor uWSGI script

description "uWSGI Emperor"
start on runlevel [2345]
stop on runlevel [06]

exec /opt/unbit/uwsgi/uwsgi /etc/uwsgi/emperor.ini
